# roles/logrotate_role/templates/logrotate-hourly.timer.j2
# ansible_managed
[Unit]
Description=Hourly rotation of log files
Documentation=man:logrotate(8) man:logrotate.conf(5)

[Timer]
OnCalendar=hourly
AccuracySec=1h
Persistent=true
Unit=logrotate.service

[Install]
WantedBy=timers.target
